# user-srv/data/etcprod/supervisor/user-srv.conf
[program:user-srv]
directory=/home/www/server/%(program_name)s
command=/home/www/server/%(program_name)s/bin/%(program_name)s --config=config/prod.toml --host=127.0.0.1
environment=EGO_LOG_PATH="/home/www/logs/user-srv",EGO_LOG_ADD_APP="true"
process_name=%(program_name)s
user=root
numprocs=1
autostart=true
startsecs=3
startretries=3
autorestart=true
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10
serverurl=AUTO
stdout_logfile=/home/system/log/supervisorlog/%(program_name)s_stdout.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
stdout_capture_maxbytes=0
stdout_events_enabled=true
stderr_logfile=/home/system/log/supervisorlog/%(program_name)s_stderr.log
stderr_logfile_maxbytes=50MB
stderr_logfile_backups=10
stderr_capture_maxbytes=0
stderr_events_enabled=false